# ansible/hellogo.service
[Unit]
Description=My Go Application

[Service]
ExecStart=/home/guobin/hellogo
Restart=always

[Install]
WantedBy=multi-user.target